# scripts/auth.cron
#! /bin/bash

/var/www/auth/auth backup
