# scripts/myaccounts.cron
#! /bin/bash

/var/www/myaccounts/myaccounts backup
